# roles/endpoint/templates/etc/systemd/system/apt-daily.timer.d/override.conf.j2
###

[Timer]
OnCalendar=
OnCalendar=*-*-* 00:01
RandomizedDelaySec=30min
AccuracySec=1m
Persistent=true
